# etc/qr.conf
[program:qe]
directory=/var/www/qr.rakuen.tokyo
command=/var/www/qr.rakuen.tokyo/qr
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/qr.log
stdout_logfile=/var/log/qr.log